# apcupsd/rootfs/etc/s6-overlay/s6-rc.d/syslogd/run
#!/command/with-contenv bashio
# ==============================================================================
# Start syslogd service
# ==============================================================================
bashio::log.info "Starting syslogd service..."
exec /sbin/syslogd -n -O - -S -t
